# DEBIAN/postinst
sudo systemctl enable trek-tracker
sudo systemctl start trek-tracker